# debian/postinst
#!/bin/sh
# postinst script for efalive
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package


case "$1" in
    configure)
        chown -R efa:efa /usr/lib/efa
        chown -R efa:efa /usr/lib/efa2
        chown -R efa:efa /home/efa/.halevt
        chown -R efa:efa /home/efa/.config
        chown efa:efa /home/efa/.xinitrc
        chmod 440 /etc/sudoers
        # update-rc.d efa start 99 2 3 4 5 . stop 01 0 1 6 .
        #aptitude -y install fontconfig-config
        #aptitude -y install plymouth
        #plymouth-set-default-theme -R efa-live
        #sed -i 's/START_DAEMON=yes/START_DAEMON=no/' /etc/default/halevt
        echo "efa     ALL=NOPASSWD: /sbin/shutdown, /opt/efalive/bin/efalive_setup.sh, /usr/sbin/dpkg-reconfigure -fgnome keyboard-configuration, /bin/date, /sbin/insserv -d ntp, /sbin/insserv -r ntp" >> /etc/sudoers
        ls -l /etc/X11
	sed -i 's/allowed_users=.*/allowed_users=anybody/' /etc/X11/Xwrapper.config
        #update-grub
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
